# deploy/systemd/followthrough-desktop-session.service
[Unit]
Description=Followthrough isolated Spark desktop
After=network-online.target

[Service]
Type=simple
Environment=DISPLAY=:99
Environment=FOLLOWTHROUGH_DESKTOP_PROFILE=/home/alhinai/browser-profiles/followthrough-desktop
ExecStart=/home/alhinai/Projects/followthrough/scripts/run-local-desktop.sh
Restart=always
RestartSec=2
TimeoutStopSec=8

[Install]
WantedBy=default.target
